# devops/systemd/mcq.service
[Unit]
Description=MCQ Battle Platform Backend Flask Service (Gunicorn)
After=network.target

[Service]
# User and Group to execute the service under (standard Ubuntu non-root user)
User=ubuntu
Group=ubuntu

# Set working directory to the cloned repository path
WorkingDirectory=/home/ubuntu/mcq-portal

# Ensure Gunicorn and dependencies run inside the virtual environment
Environment="PATH=/home/ubuntu/mcq-portal/venv/bin"
# Load environment variables from the local .env configuration file
EnvironmentFile=/home/ubuntu/mcq-portal/.env

# Gunicorn execution command
# -w 3: Runs 3 worker processes (Formula: 2 * CPU_cores + 1)
# -b 127.0.0.1:5000: Binds to port 5000 locally (accessible by Nginx)
# run:app: Refers to run.py and the instantiated 'app' object inside it
ExecStart=/home/ubuntu/mcq-portal/venv/bin/gunicorn -w 3 -b 127.0.0.1:5000 run:app

# Auto-Restart configurations
Restart=always
RestartSec=5

# Logging configuration (directing output to journald/syslog)
StandardOutput=journal
StandardError=journal
SyslogIdentifier=mcq-backend

# Security/Resource configurations
LimitNOFILE=65536
PrivateTmp=true

[Install]
WantedBy=multi-user.target
